# deploy/cruza-warm-cache.service
[Unit]
Description=Cruza Dados - Cache Warmer
After=postgresql.service network.target
Wants=postgresql.service

[Service]
Type=simple
User=govbr
WorkingDirectory=/home/govbr/govbr-project
EnvironmentFile=/home/govbr/govbr-project/.env
ExecStart=/home/govbr/govbr-project/venv/bin/python -m web.warm_cache --daemon --loop
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
